# install/ProjectorAutomation/Library/LaunchAgents/no.ntnu.ime.blestvegg-on.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>no.ntnu.ime.blestvegg-on</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/nummedal/blestvegg.py</string>
		<string>power on</string>
	</array>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>2</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>3</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>4</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>5</integer>
			<key>Hour</key>
			<integer>7</integer>
			<key>Minute</key>
			<integer>45</integer>
		</dict>
	</array>
	<key>ExitTimeOut</key>
	<integer>600</integer>
	<key>StandardOutPath</key>
	<string>/Users/show/blestvegg/blestvegg.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/show/blestvegg/blestvegg-err.log</string>
</dict>
</plist>
